# sweetadb/systemd/boarnet-sweetadb.service
[Unit]
Description=BoarNet sweetADB adapter — ships events.jsonl into BoarNet
Documentation=https://github.com/Bino97/boarnet-adapters/tree/main/sweetadb
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# DynamicUser + StateDirectory — no pre-creation required.
# systemd allocates a transient UID at service start, persists the
# same UID across restarts for this service, and creates
# /var/lib/boarnet-sweetadb owned by it (mode 0700). Lifetime of
# the user tracks the lifetime of the service. Eliminates the
# status=217/USER failure mode that bit early partners who
# deployed before creating a `boarnet` system user.
DynamicUser=true
StateDirectory=boarnet-sweetadb
StateDirectoryMode=0700

# Default log paths — override in /etc/boarnet/sweetadb.env (or a
# systemd drop-in, see README) if sweetADB writes elsewhere. Values
# here are defaults; EnvironmentFile= below takes precedence.
Environment=SWEETADB_EVENTS_LOG=/var/log/sweetadb/mimic/events.jsonl
Environment=SWEETADB_PAYLOADS_DIR=/var/log/sweetadb/mimic/payloads

EnvironmentFile=/etc/boarnet/sweetadb.env
ExecStart=/usr/local/bin/boarnet-sweetadb-adapter \
  --events-log ${SWEETADB_EVENTS_LOG} \
  --payloads-dir ${SWEETADB_PAYLOADS_DIR} \
  --data-dir /var/lib/boarnet-sweetadb \
  --sensor-id ${BOARNET_SENSOR_ID} \
  --token ${BOARNET_TOKEN}
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Hardening — the adapter needs read on the honeypot log + write
# on its StateDirectory + outbound HTTPS. Nothing else.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Read access to the default log location. If you point
# SWEETADB_EVENTS_LOG at a path outside /var/log/sweetadb (e.g.
# /root/mimic/), you ALSO need a drop-in to grant read access —
# see the README "Non-default log locations" section.
ReadOnlyPaths=/var/log/sweetadb

[Install]
WantedBy=multi-user.target
